# opt/setup/templates/systemd/adblock-sync.service
[Unit]
Description=Apply Adblock nftables rules
After=network-online.target ssh-adblock-dns.service
Wants=network-online.target ssh-adblock-dns.service
ConditionPathExists=__SSH_DNS_ADBLOCK_SYNC_BIN__

[Service]
Type=oneshot
ExecStart=__SSH_DNS_ADBLOCK_SYNC_BIN__ --apply

[Install]
WantedBy=multi-user.target
