# opt/setup/templates/systemd/openvpn-speed.service
[Unit]
Description=OpenVPN per-user speed limiter (tc)
After=network-online.target __OVPN_TCP_SERVICE_NAME__
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/openvpn-speed watch --env-file __OPENVPN_RUNTIME_ENV_FILE__ --interval __OVPN_SPEED_INTERVAL__
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
