# opt/setup/templates/systemd/openvpn-speed-reconcile.path
[Unit]
Description=Watch OpenVPN speed inputs for immediate reconcile

[Path]
PathChanged=__OPENVPN_SPEED_EVENT_DIR__
PathChanged=__OPENVPN_SPEED_STATE_ROOT__
Unit=__OPENVPN_SPEED_RECONCILE_SERVICE__.service

[Install]
WantedBy=multi-user.target
