# opt/setup/templates/systemd/ovpn-ws-proxy.service
[Unit]
Description=OpenVPN Websocket Proxy (Go) service
After=network-online.target __OPENVPN_TCP_SERVICE__.service
Requires=__OPENVPN_TCP_SERVICE__.service
PartOf=__OPENVPN_TCP_SERVICE__.service
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-__OPENVPN_CONFIG_ENV_FILE__
ExecStart=/usr/local/bin/ws-proxy --mode tcp --listen-host 127.0.0.1 --listen-port ${OPENVPN_WS_PROXY_PORT} --backend-host 127.0.0.1 --backend-port ${OPENVPN_PORT_TCP} --path ${OPENVPN_WS_PUBLIC_PATH} --handshake-timeout ${OPENVPN_WS_HANDSHAKE_TIMEOUT_SEC}
Restart=always
RestartSec=2
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
