# opt/setup/templates/systemd/ssh-expired-cleaner.service
[Unit]
Description=Delete expired SSH-linked accounts and linked OpenVPN artifacts
After=network-online.target bot-telegram-backend.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ssh-expired-cleaner --once
User=root
Group=root
Nice=10
