# opt/setup/templates/systemd/sshws-qac-enforcer.service
[Unit]
Description=SSH WS quota and access enforcer
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/sshws-qac-enforcer --once
User=root
Group=root
Nice=10
