# opt/setup/templates/systemd/sshws-stunnel.service
[Unit]
Description=Stunnel local TLS bridge for SSH WebSocket
After=sshws-dropbear.service
Requires=sshws-dropbear.service

[Service]
Type=simple
ExecStartPre=/usr/bin/mkdir -p /run/stunnel
ExecStart=__STUNNEL_BIN__ /etc/stunnel/sshws.conf
Restart=always
RestartSec=2
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
