# opt/setup/templates/systemd/warp-zt-socks-bridge.service
[Unit]
Description=WARP Zero Trust SOCKS bridge for Xray
After=network-online.target __WARP_ZEROTRUST_SERVICE__
Wants=network-online.target

[Service]
Type=simple
Environment=WARP_ZT_BRIDGE_LISTEN_HOST=127.0.0.1
Environment=WARP_ZT_BRIDGE_LISTEN_PORT=__WARP_ZEROTRUST_BRIDGE_PORT__
Environment=WARP_ZT_BRIDGE_UPSTREAM_HOST=127.0.0.1
Environment=WARP_ZT_BRIDGE_UPSTREAM_PORT=__WARP_PROXY_PORT__
ExecStart=__WARP_ZEROTRUST_BRIDGE_BIN__
Restart=always
RestartSec=2
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
