# opt/setup/templates/systemd/xray-quota.service
[Unit]
Description=Xray quota watcher (metadata -> auto block)
After=network-online.target xray.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/xray-quota watch --config /usr/local/etc/xray/conf.d/30-routing.json --marker dummy-quota-user --interval 2
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
